# templates/tag-notify.service.erb
[Unit]
Description=Watch facter file with inotify and execute Puppet run on change
ConditionPathExists=<%= @fact_file %>

[Service]
Type=simple
ExecStart=<%= @inotify_script_path %> "<%= @fact_file %>"
ProtectHome=true
ProtectKernelTunables=true
ProtectSystem=full
RemainAfterExit=no

[Install]
WantedBy=multi-user.target